# build/worker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:queue]
command=php /var/www/artisan queue:work --sleep=3 --tries=3 --timeout=300
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:schedule]
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

; enable only if you use horizon
; [program:horizon]
; command=php /var/www/artisan horizon
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr

; enable only if you use reverb
; [program:reverb]
; command=php /var/www/artisan reverb:start --host=0.0.0.0 --port=8080
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr