# build/worker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

; Horizon supervises the redis queue workers (M4); plain queue:work is disabled.
; [program:queue]
; command=php /var/www/artisan queue:work --sleep=3 --tries=3 --timeout=300
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr

[program:schedule]
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:horizon]
command=php /var/www/artisan horizon
autostart=true
autorestart=true
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; enable only if you use reverb
; [program:reverb]
; command=php /var/www/artisan reverb:start --host=0.0.0.0 --port=8080
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr